# deploy/promptgrimoire-worker.service
[Unit]
Description=PromptGrimoire Export Worker
After=network-online.target postgresql.service pgbouncer.service promptgrimoire.service
Wants=network-online.target
Requires=postgresql.service
PartOf=promptgrimoire.service

[Service]
Type=notify
User=promptgrimoire
Group=promptgrimoire
WorkingDirectory=/opt/promptgrimoire
EnvironmentFile=/opt/promptgrimoire/.env
Environment=PATH=/home/promptgrimoire/.local/bin:/home/promptgrimoire/.TinyTeX/bin/x86_64-linux:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/promptgrimoire
Environment=OSFONTDIR=/usr/share/fonts:/usr/share/texmf/fonts

ExecStart=/home/promptgrimoire/.local/bin/uv run python -m promptgrimoire.export.worker_main

# Watchdog: restart if no heartbeat within 5 minutes
WatchdogSec=300
Restart=always
# Covers watchdog timeouts, crashes, AND PartOf-triggered stops
# (when the app service restarts, PartOf stops the worker cleanly;
# Restart=always brings it back with fresh code)
RestartSec=5

# Shutdown: worker cancels in-flight job on SIGTERM. TimeoutStopSec gives
# headroom for asyncio cleanup and DB connection teardown.
TimeoutStopSec=30
KillSignal=SIGTERM

# Resource controls: best-effort scheduling
Nice=19
IOSchedulingClass=idle
CPUWeight=10
MemoryMax=3G
MemoryHigh=2560M
OOMScoreAdjust=500

# Security hardening
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/promptgrimoire/logs /opt/promptgrimoire/.venv /home/promptgrimoire/.TinyTeX /home/promptgrimoire/.cache/uv /tmp
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
